# udev/61-cdplayer-start.rules
# Trigger play audio script on CD insertion
ACTION=="change", SUBSYSTEM=="block", KERNEL=="sr*", ID_CDROM_MEDIA_CD=="1", RUN+="/lib/udev/scripts/startplaycdservice.sh"
